# go-gateway/deployment/edge/systemd/bifrost-gateway.service
[Unit]
Description=Bifrost Industrial Gateway
Documentation=https://docs.bifrost.dev
After=network.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=bifrost
Group=bifrost
ExecStart=/opt/bifrost/bin/bifrost-gateway -config /etc/bifrost/gateway.yaml
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
TimeoutStopSec=30

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/bifrost /var/lib/bifrost /tmp
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
RestrictNamespaces=true
SystemCallArchitectures=native

# Environment
Environment=BIFROST_CONFIG_PATH=/etc/bifrost/gateway.yaml
Environment=BIFROST_LOG_LEVEL=info
Environment=BIFROST_DATA_DIR=/var/lib/bifrost
Environment=BIFROST_LOG_DIR=/var/log/bifrost

# Working directory
WorkingDirectory=/opt/bifrost

# Standard output and error
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bifrost-gateway

[Install]
WantedBy=multi-user.target